# notice.ini
[program:notice]
directory = /usr/local/services/notice/
command = gunicorn -w 4 -b 0.0.0.0:5000 manager:app
autostart = true
autorestart = true
startsecs = 10
startretries = 10
stopasgroup = true
killasgroup = true
redirect_stderr = true
stdout_logfile = /usr/local/services/app-log/notice.log
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 20